# zpool-scrub@.timer
[Unit]
Description=Scrub ZFS pool weekly

[Timer]
OnCalendar=*-*-1,15 01:00:00
Persistent=true

[Install]
WantedBy=timers.target
